# debian/varnishncsa.service
[Unit]
Description=Varnish HTTP accelerator log daemon
After=varnish.service

[Service]
User=varnishlog
ExecStart=/usr/bin/varnishncsa -a -w /var/log/varnish/varnishncsa.log
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
